# systemd/simplemon.service
[Unit]
Description=SimpleMon – lightweight system monitoring daemon
Documentation=https://github.com/simplemon
After=network.target
Wants=network.target

[Service]
Type=simple
User=simplemon
Group=simplemon

ExecStart=/usr/local/bin/simplemon -config /etc/simplemon/simplemon.yaml
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=simplemon

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/simplemon
CapabilityBoundingSet=

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
